# config/crontab/checkOldOrder.sh
#!/bin/bash
curl 'http://localhost/orders/api/checkOldOrder'
